# filesystem-backup/usr/lib/opkg/info/mtkhqos_util.list
/etc/config/mtkhqos
/etc/init.d/mtkhqos
/sbin/mtkhqos
