# filesystem-backup/usr/lib/opkg/info/switch.list
/usr/sbin/switch_ctl.sh
/etc/init.d/switch_ctl
/usr/sbin/switch
